# automation/com.health-analytics.daily-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.health-analytics.daily-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/nick/clawd/projects/health-analytics/venv/bin/python</string>
        <string>/Users/nick/clawd/projects/health-analytics/health</string>
        <string>generate</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/nick/clawd/projects/health-analytics</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/nick/clawd/projects/health-analytics/logs/daily-refresh.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/nick/clawd/projects/health-analytics/logs/daily-refresh.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/nick</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
